# vuln-monitor.timer
[Unit]
Description=Run vuln-monitor every 5 minutes
Requires=vuln-monitor.service

[Timer]
# fire 2 min after boot, then every 5 min after the last run completes
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
# if the machine was off when a fire was due, run once on boot
Persistent=true
Unit=vuln-monitor.service

[Install]
WantedBy=timers.target
