# systemd/vuln-monitor.service
[Unit]
Description=Vulnerability intelligence aggregator (oneshot)
Documentation=https://github.com/Knaithe/1DayNews
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=vuln
Group=vuln
WorkingDirectory=/opt/vuln-monitor
EnvironmentFile=/opt/vuln-monitor/.env
Environment=VULN_DATA_DIR=/opt/vuln-monitor
ExecStart=/opt/vuln-monitor/venv/bin/python /opt/vuln-monitor/src/vuln_monitor.py fetch --no-push
ExecStart=/opt/vuln-monitor/venv/bin/python /opt/vuln-monitor/src/vuln_monitor.py enrich

# resource / time caps (increased for LLM enrichment)
RuntimeMaxSec=900
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# sandboxing
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/vuln-monitor
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
