# systemd/vuln-web.service
[Unit]
Description=vuln-monitor web dashboard
After=network.target

[Service]
Type=simple
User=vuln
Group=vuln
WorkingDirectory=/opt/vuln-monitor
Environment=VULN_DATA_DIR=/opt/vuln-monitor
ExecStart=/opt/vuln-monitor/venv/bin/python /opt/vuln-monitor/src/web.py --public
Restart=on-failure
RestartSec=5

# sandboxing — DB stays read-only, but SQLite WAL mode needs SHM/WAL files writable
# (`-` prefix = ignore if missing, in case a fresh DB has not generated them yet)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vuln-monitor/.web_token -/opt/vuln-monitor/vuln_cache.db-shm -/opt/vuln-monitor/vuln_cache.db-wal
ReadOnlyPaths=/opt/vuln-monitor
PrivateTmp=true

[Install]
WantedBy=multi-user.target
