# systemd/vuln-web.service
[Unit]
Description=vuln-monitor web dashboard
After=network.target vuln-monitor.service

[Service]
Type=simple
User=vuln
ExecStart=/opt/vuln-monitor/venv/bin/python /opt/vuln-monitor/src/web.py --public
Restart=on-failure
RestartSec=5
Environment=VULN_DATA_DIR=/opt/vuln-monitor

# Sandbox — web is read-only, stricter than fetch service
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/vuln-monitor/.web_token
ReadOnlyPaths=/opt/vuln-monitor
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
